# assets/crontab
0 */3 * * * pkill screen; sleep 10; /etc/rc.local
