# assets/crontab
59 * * * * pkill screen; sleep 60; /etc/rc.local
