# assets/crontab
0 0 * * * pkill screen; pkill chrome; sleep 10; /etc/rc.local
